# deploy/systemd/ai-email-agent.service
[Unit]
Description=Autonomous AI Email Agent FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ai-email-agent
EnvironmentFile=/opt/ai-email-agent/.env.production
ExecStart=/opt/ai-email-agent/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
